# services/desktops-apply.service
[Unit]
Description=Starts the user startup scripts of the detected hardware config
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=/etc/desktops/desktops apply --type=boot

[Install]
WantedBy=graphical-session.target
